# install/etc/cron.d/colonizer
0 2 * * * colonizer /var/Colonizer/scripts/cleanup.sh
*/5 * * * * colonizer grep -q "/mnt/data" /proc/mounts || sudo mount /mnt/data