# install/etc/supervisor/conf.d/nginx.conf
[program:nginx]
command=/usr/sbin/nginx
autostart=true
stdout_logfile=/app/Colonizer/log/nginx.out.log
stderr_logfile=/app/Colonizer/log/nginx.err.log