# install/etc/supervisor/conf.d/nginx.conf
[program:nginx]
command=/usr/sbin/nginx
autostart=true
stdout_logfile=/var/log/colonizer-nginx.out.log
stderr_logfile=/var/log/colonizer-nginx.err.log